# infra/docker/s6/s6-rc.d/control-plane/finish
#!/bin/sh
# Allow time for graceful shutdown
sleep 5
