# infra/docker/s6/s6-rc.d/control-plane/run
#!/command/with-contenv sh

# Wait for PostgreSQL only if using embedded (local Unix socket) mode
if [ ! -f /tmp/skip-embedded-postgres ]; then
  while ! gosu postgres pg_isready -h /var/run/postgresql -q 2>/dev/null; do
    sleep 1
  done
  export DATABASE_SOCKET=/var/run/postgresql
fi

# Wait for Redis only if using embedded (local Unix socket) mode
if [ ! -f /tmp/skip-embedded-redis ]; then
  while ! redis-cli -s /var/run/redis/redis.sock -a "${REDIS_PASSWORD:-honorclaw}" ping 2>/dev/null | grep -q PONG; do
    sleep 1
  done
  export REDIS_SOCKET=/var/run/redis/redis.sock
  export REDIS_PASSWORD="${REDIS_PASSWORD:-honorclaw}"
fi

export NODE_ENV=${NODE_ENV:-production}

# Drop privileges: run the Node.js app as the 'node' user (created by node:alpine base image)
exec gosu node node /app/packages/control-plane/dist/main.js
