# infra/docker/s6/s6-rc.d/control-plane/run
#!/bin/sh

# Wait for PostgreSQL only if using embedded (local Unix socket) mode
if [ ! -f /tmp/skip-embedded-postgres ]; then
  while ! su-exec postgres pg_isready -h /var/run/postgresql -q 2>/dev/null; do
    sleep 1
  done
  export DATABASE_SOCKET=/var/run/postgresql
fi

# Wait for Redis only if using embedded (local Unix socket) mode
if [ ! -f /tmp/skip-embedded-redis ]; then
  while ! redis-cli -s /var/run/redis/redis.sock -a "${REDIS_PASSWORD:-honorclaw}" ping 2>/dev/null | grep -q PONG; do
    sleep 1
  done
  export REDIS_SOCKET=/var/run/redis/redis.sock
fi

export NODE_ENV=${NODE_ENV:-production}

exec node /app/packages/control-plane/dist/main.js
