# infra/docker/s6/s6-rc.d/ollama/run
#!/bin/sh
if [ "${OLLAMA_DISABLED:-false}" = "true" ]; then
  exec sleep infinity
fi

mkdir -p /data/ollama/models
export OLLAMA_MODELS=/data/ollama/models

# Start Ollama in the background, then pull the default model
ollama serve &
OLLAMA_PID=$!

# Wait for Ollama to become ready (up to 30 seconds)
OLLAMA_BASE_URL="${OLLAMA_BASE_URL:-http://localhost:11434}"
for i in $(seq 1 30); do
  if curl -sf "${OLLAMA_BASE_URL}/api/tags" >/dev/null 2>&1; then
    break
  fi
  sleep 1
done

# Discover the latest Llama model if no explicit default is set.
# Probes the Ollama registry for the newest version and falls back gracefully.
if [ -z "${HONORCLAW_DEFAULT_MODEL}" ]; then
  REGISTRY="https://registry.ollama.ai/v2/library"
  DISCOVERED=""
  for candidate in llama4 llama3.3 llama3.2 llama3.1 llama3; do
    if curl -sf "${REGISTRY}/${candidate}/tags/list" >/dev/null 2>&1; then
      DISCOVERED="${candidate}"
      break
    fi
  done
  DEFAULT_MODEL="${DISCOVERED:-llama3.2}"
  echo "Discovered latest Llama model: ${DEFAULT_MODEL}"
else
  DEFAULT_MODEL="${HONORCLAW_DEFAULT_MODEL}"
fi

# Pull the default model if not already present
if ! ollama list 2>/dev/null | grep -q "^${DEFAULT_MODEL}"; then
  echo "Pulling default model: ${DEFAULT_MODEL} ..."
  ollama pull "${DEFAULT_MODEL}" 2>&1 || echo "WARNING: Failed to pull ${DEFAULT_MODEL}"
fi

# Persist the resolved model name so the control plane can read it
echo "${DEFAULT_MODEL}" > /data/ollama/default-model

# Keep Ollama running as the foreground process
wait $OLLAMA_PID
